# linux/common/theme-watch/.config/systemd/user/generate-theme.path
[Unit]
Description=Watch theme/palette.toml and regenerate the theme on change

[Path]
PathModified=%h/dotfiles/theme/palette.toml
Unit=generate-theme.service

[Install]
WantedBy=default.target
