# linux/common/theme-watch/.config/systemd/user/generate-theme.service
[Unit]
Description=Regenerate terminal theme from theme/palette.toml

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/generate-theme.py
